# rules.d/91-hsdpa-modem.rules
ACTION=="add", SUBSYSTEM=="tty", KERNEL=="ttyUSB0", RUN+="/usr/local/bin/umtson"
