# rules.d/92-usb-modem.rules
SUBSYSTEM=="usb", ACTION=="add", SYSFS{idVendor}=="12d1", SYSFS{idProduct}=="1446", RUN+="/usr/sbin/usb_modeswitch"
